# tapo-camera-mcp-watchfiles.service
[Unit]
Description=Tapo Camera MCP WebApp with Watchfiles Crashproofing
After=network.target
Wants=network.target

[Service]
Type=simple
User=tapo-camera-mcp
Group=tapo-camera-mcp
WorkingDirectory=/opt/tapo-camera-mcp
Environment=PATH=/opt/tapo-camera-mcp/venv/bin
Environment=TAPO_WEBAPP_HOST=0.0.0.0
Environment=TAPO_WEBAPP_PORT=7777
Environment=TAPO_WEBAPP_DEBUG=false
Environment=WATCHFILES_MAX_RESTARTS=20
Environment=WATCHFILES_RESTART_DELAY=2.0
Environment=WATCHFILES_BACKOFF_MULTIPLIER=1.5
Environment=WATCHFILES_HEALTH_CHECK_INTERVAL=60
Environment=WATCHFILES_NOTIFY_ON_CRASH=true
ExecStart=/opt/tapo-camera-mcp/venv/bin/python watchfiles_runner.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tapo-camera-mcp-watchfiles

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/tapo-camera-mcp/logs /opt/tapo-camera-mcp/data
ProtectHome=yes

# Resource limits
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
